# dotfiles/.config/autostart/mega-sync.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Exec=/run/current-system/sw/bin/megasync %u
